# readme/etc/systemd/system/SplunkCopyFiles.service
[Unit]
Description="Move Splunk Filesystem Output"

[Service]
ExecStart=/home/splunk/move_files.sh
User=splunk
Group=splunk
